# conf/peanut_supervisor.conf
[program:peanutapp]
command = uwsgi --ini /var/www/peanut/conf/peanut_uwsgi.ini

[program:nginx]
command = service nginx restart
